# ccd/ansible/erikube/roles/kube-master/files/etc/profile.d/kubeconfig.sh
# kubectl env variable

export KUBECONFIG=/etc/kubernetes/admin.conf
